# ansible/roles/borg/templates/borg-backup.timer.j2
[Unit]
Description=Run borg backup

[Timer]
OnCalendar=*-*-* {{ borg_backup_run_hour }}:00
Persistent=true

[Install]
WantedBy=timers.target
