# ansible/roles/sanoid/templates/syncoid.service.j2
[Unit]
Description=Run syncoid
After=network.target
[Service]
Type=oneshot
User={{ sanoid_syncoid_user }}
Group={{ sanoid_syncoid_user }}
ExecStart=/opt/syncoid/syncoid.sh
