# ansible/roles/borg/borg-backup.service.j2
[Unit]
Description=Run borg backup
After=network.target
[Service]
Type=oneshot
ExecStart=/opt/borg/backup.sh
